# services/shoplist.service
[Unit]
Description=shoplist
After=network-online.target
Wants=network-online.target

[Service]
User=shoplist
Group=shoplist
EnvironmentFile=/etc/shoplist.env
ExecStart=/usr/local/bin/shoplist
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/shoplist
WorkingDirectory=/var/lib/shoplist

[Install]
WantedBy=multi-user.target
